# meta-local/recipes-graphics/kiosker/kiosker/kiosker.socket
[Socket]
ListenDatagram=/run/kiosker.sock
Accept=no

[Install]
WantedBy=sockets.target

